# example/wikiwash-fetch.cron
0 * * * * deploy /bin/bash -l -c "cd /home/deploy/wikiwash/current && /usr/local/bin/node scripts/FetchTopArticles.js"
